# template_mailcow/7.4/files/mailcow-monitor.timer
[Unit]
Description=Mailcow Monitoring Collector Timer

[Timer]
OnBootSec=30
OnUnitActiveSec=60
AccuracySec=5

[Install]
WantedBy=timers.target
