# Operating Systems/Linux/template_zypper_updateinfo/6.0/files/systemd/zabbix-template-module-zypper-updateinfo.timer
[Unit]
Description=Send Zypper updateinfo to Zabbix
Requires=zabbix-template-module-zypper-updateinfo.service

[Timer]
Unit=zabbix-template-module-zypper-updateinfo.service
OnCalendar=*-*-* *:00:00

[Install]
WantedBy=timers.target