# Applications/Clustered_File_Systems/Ceph/template_ceph_cache_tier_sender/7.4/files/ceph-osd-perf-zabbix@.service
[Unit]
Description=Ceph OSD metrics to Zabbix (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/flock -n /run/ceph-osd-perf-zabbix-%i.lock /usr/local/sbin/ceph-osd-perf-zabbix --%i
