# Applications/Clustered_File_Systems/Ceph/template_ceph_cache_tier_sender/7.4/files/ceph-osd-perf-zabbix-perf-raw.timer
[Unit]
Description=Ceph OSD raw perf dump to Zabbix (optional, every 60 seconds)

[Timer]
OnBootSec=45s
OnUnitActiveSec=60s
AccuracySec=5s
Unit=ceph-osd-perf-zabbix@perf-raw.service

[Install]
WantedBy=timers.target
